# .config/systemd/user/aw-qt.service
[Unit]
Description=ActivityWatch

[Service]
ExecStart=/usr/bin/aw-qt

[Install]
WantedBy=wm-session.target
