# .config/systemd/user/blueman.service
[Unit]
Description=Blueman

[Service]
ExecStart=/usr/bin/blueman-applet

[Install]
WantedBy=wm-session.target
