# .config/systemd/user/flameshot.service
[Unit]
Description=Flameshot

[Service]
ExecStart=/usr/bin/flameshot

[Install]
WantedBy=wm-session.target
