# .config/systemd/user/nm-applet.service
[Unit]
Description=Network Manager Applet

[Service]
ExecStart=/usr/bin/nm-applet

[Install]
WantedBy=wm-session.target
