# .config/systemd/user/notmuch.timer
[Unit]
Description=Run notmuch sync script every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
