# .config/systemd/user/notmuch.timer
[Unit]
Description=Run notmuch new every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
