# .config/systemd/user/opensnitch-ui.service
[Unit]
Description=OpenSnitch UI

[Service]
ExecStart=/usr/bin/opensnitch-ui

[Install]
WantedBy=wm-session.target
