# .config/systemd/user/snixembed.service
[Unit]
Description=snixembed

[Service]
ExecStart=/usr/bin/snixembed
Restart=on-failure

[Install]
WantedBy=wm-session.target
