# .config/systemd/user/vnstat.service
[Unit]
Description=vnStat network monitor

[Service]
ExecStart=/usr/bin/vnstatd

[Install]
WantedBy=wm-session.target
