# .config/systemd/user/wm-session.target
[Unit]
Description=Window Manager Session
After=graphical-session.target
