# .config/systemd/user/xsettingsd.service
[Unit]
Description=X Settings Daemon

[Service]
ExecStart=/usr/bin/xsettingsd

[Install]
WantedBy=wm-session.target
